# dot_config/systemd/user/lid-monitor.service.tmpl
[Unit]
Description=Laptop Lid Monitor for Display Management
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.config/hypr/scripts/lid_monitor.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical-session.target
